# pkg/payload/Library/LaunchAgents/io.github.xishie.helios.InternalNetworkAvailable.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>KeepAlive</key>
        <true/>
        <key>Label</key>
        <string>io.github.xishie.helios.InternalNetworkAvailable</string>
        <key>RunAtLoad</key>
        <true/>
        <key>ProgramArguments</key>
        <array>
                <string>/Library/helios/sso_bundle.app/Contents/MacOS/sso_bundle</string>
                <string>-notification</string>
                <string>com.apple.KerberosPlugin.InternalNetworkAvailable</string>
                <string>-action</string>
                <string>/Library/helios/helios.sh</string>
        </array>
</dict>
</plist>