# pkg/payload/Library/LaunchAgents/io.github.xishie.helios.timer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>io.github.xishie.helios.timer</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Library/helios/helios.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/tmp/io.github.xishie.helios.timer.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/io.github.xishie.helios.timer.err</string>
  </dict>
</plist>
